# backend/DatabaseSync/cron.crontab
0 * * * * /usr/bin/node /path/to/your/databasePull.js
15 * * * * /usr/bin/node /path/to/your/ipfsBatchDownloader.js
30 * * * * /usr/bin/python /path/to/your/updateMongo.py
